# build/faas-manager/supervisord.conf
[supervisord]
nodaemon=true

[program:faas-manager]
command=/root/.local/bin/faas-manager
autorestart=true
startretries=50
priority=0
